# target/maven-status/jdbcdemo.service
[Unit]
Description=myblog
After=syslog.target

[Service]
ExecStart= /usr/bin/java -jar/root/demo/jdbcdemo.jar

[Install]
WantedBy=multi-user.target
